# debian/prerm
#!/bin/bash
set -e

# Pre-removal script for pepctl
# This script runs before the package is removed

case "$1" in
    remove|upgrade|deconfigure)
        # Stop pepctl service if it's running
        if systemctl is-active --quiet pepctl 2>/dev/null; then
            echo "Stopping pepctl service..."
            systemctl stop pepctl || true
        fi
        
        # Disable pepctl service if we're removing (not upgrading)
        if [ "$1" = "remove" ]; then
            if systemctl is-enabled --quiet pepctl 2>/dev/null; then
                echo "Disabling pepctl service..."
                systemctl disable pepctl || true
            fi
        fi
        ;;
        
    failed-upgrade)
        # Nothing to do
        ;;
        
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0 